# deploy/arbitrage-altyn-bot.service
[Unit]
Description=Altyn/Rapira Telegram spread bot
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
User=arbitrage-bot
Group=arbitrage-bot
WorkingDirectory=/opt/arbitrage-altyn-bot
EnvironmentFile=/etc/arbitrage-altyn-bot/bot.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/arbitrage-altyn-bot/.venv/bin/arbitrage-altyn-bot
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
UMask=0077
StateDirectory=arbitrage-altyn-bot
StateDirectoryMode=0700
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arbitrage-altyn-bot
TasksMax=64

CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectClock=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
ReadWritePaths=/var/lib/arbitrage-altyn-bot
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
